# config/cron.d/pdf
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

*/20 * * * * root /usr/local/bin/pdf-mpc-diff-to-n8n >/dev/null 2>&1
0 */6 * * * root /usr/local/bin/pdf-jpl-diff-to-n8n >/dev/null 2>&1
*/30 * * * * root /usr/local/bin/pdf-heartbeat >/dev/null 2>&1
