# infra/proxmox-export-2026-03-14/ct124/pdf-orbital-diagram-build.service
[Unit]
Description=Build PDF live orbital diagram
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
ExecStart=/opt/pdf/orbital-diagram/build_orbital_diagram.py
