# infra/proxmox-export-2026-03-14/ct124/pdf-trajectory-baseline.service
[Unit]
Description=PDF Trajectory Baseline Job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/pdf/trajectory/.venv/bin/python /opt/pdf/trajectory/trajectory_baseline.py
User=root
