# infra/proxmox-export-2026-03-14/ct124/pdf-orbital-advanced.service
[Unit]
Description=Build PDF advanced orbital visualization
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
ExecStart=/opt/pdf/orbital-diagram/advanced_orbit.py
