# infra/proxmox-export-2026-03-14/host/pdf-dr-drill.service
[Unit]
Description=PDF DR Drill Report
[Service]
Type=oneshot
ExecStart=/root/pdf-drill/run_drill.sh
