# infra/proxmox-export-2026-03-14/ct123/pdf-briefing-bot.service
[Unit]
Description=PDF AI Briefing Bot
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
EnvironmentFile=/etc/pdf-alert-brain.env
ExecStart=/opt/pdf/alert-brain/briefing_bot.py
