# infra/proxmox-export-2026-03-14/ct123/pdf-briefing-bot.timer
[Unit]
Description=Run PDF AI Briefing every hour
[Timer]
OnBootSec=6min
OnUnitActiveSec=1h
Persistent=true
[Install]
WantedBy=timers.target
