# infra/proxmox-export-2026-03-14/ct123/pdf-threshold-tuning.timer
[Unit]
Description=Run PDF Threshold Tuning every 30 min
[Timer]
OnBootSec=8min
OnUnitActiveSec=30min
Persistent=true
[Install]
WantedBy=timers.target
