# infra/proxmox-export-2026-03-14/ct123/pdf-policy-router.timer
[Unit]
Description=Run PDF Policy Router every 15 min
[Timer]
OnBootSec=7min
OnUnitActiveSec=15min
Persistent=true
[Install]
WantedBy=timers.target
