# infra/proxmox-export-2026-03-14/ct123/pdf-kpi-board.service
[Unit]
Description=PDF KPI Board Build
[Service]
Type=oneshot
ExecStart=/opt/pdf/expansion/kpi_board.py
