# infra/proxmox-export-2026-03-14/ct123/pdf-threshold-tuning.service
[Unit]
Description=PDF Threshold Tuning Lab
[Service]
Type=oneshot
ExecStart=/opt/pdf/expansion/threshold_tuning.py
