# infra/proxmox-export-2026-03-14/ct123/pdf-incident-generator.timer
[Unit]
Description=Run PDF Incident Generator every 15 min
[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true
[Install]
WantedBy=timers.target
