# infra/proxmox-export-2026-03-14/ct123/pdf-token-rotation.service
[Unit]
Description=PDF Monthly Token Rotation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/pdf/expansion/rotate_tokens.sh
