# infra/proxmox-export-2026-03-14/ct123/pdf-executive-brief.service
[Unit]
Description=PDF Executive Brief Generator
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
ExecStart=/opt/pdf/next3/executive_brief.py
