# infra/proxmox-export-2026-03-14/ct123/pdf-export-feeds.service
[Unit]
Description=PDF Partner Export Feeds
[Service]
Type=oneshot
ExecStart=/opt/pdf/expansion/export_feeds.py
