# infra/proxmox-export-2026-03-14/ct123/pdf-night-watch.timer
[Unit]
Description=Run PDF Night Watch every 10 min
[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true
[Install]
WantedBy=timers.target
