# infra/proxmox-export-2026-03-14/ct123/pdf-integrity-chain.timer
[Unit]
Description=Run PDF Integrity Chain every 15 min
[Timer]
OnBootSec=9min
OnUnitActiveSec=15min
Persistent=true
[Install]
WantedBy=timers.target
