# infra/proxmox-export-2026-03-14/ct123/pdf-sre-monitor.service
[Unit]
Description=PDF SRE monitor snapshot
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
ExecStart=/opt/pdf/sre-monitor/sre_monitor.py
