# infra/proxmox-export-2026-03-14/ct123/pdf-alert-brain.timer
[Unit]
Description=Run PDF Alert Brain every 15 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
