# infra/proxmox-export-2026-03-14/ct123/pdf-trend-engine.service
[Unit]
Description=PDF Trend Engine
After=pdf-historian-snapshot.service
Wants=pdf-historian-snapshot.service
[Service]
Type=oneshot
ExecStart=/opt/pdf/mission-control/trend_engine.py
