# infra/proxmox-export-2026-03-14/ct123/pdf-incident-generator.service
[Unit]
Description=PDF Incident Generator
After=pdf-alert-brain.service
Wants=pdf-alert-brain.service
[Service]
Type=oneshot
ExecStart=/opt/pdf/mission-control/incident_generator.py
