# infra/proxmox-export-2026-03-14/ct123/pdf-policy-router.service
[Unit]
Description=PDF Notification Policy Router
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
Environment=WEBHOOK_CRITICAL=REDACTED_WEBHOOK_URL
Environment=WEBHOOK_WATCH=REDACTED_WEBHOOK_URL
ExecStart=/opt/pdf/expansion/policy_router.py
