# infra/proxmox-export-2026-03-14/ct123/pdf-mission-control.service
[Unit]
Description=PDF Mission Control API
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/pdf-secrets.env
Type=simple
WorkingDirectory=/opt/pdf/mission-control
ExecStart=/opt/pdf/alert-brain/.venv/bin/uvicorn mission_control_api:app --host 0.0.0.0 --port 9100
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
