# infra/proxmox-export-2026-03-14/ct123/pdf-integrity-chain.service
[Unit]
Description=PDF Integrity Chain Update
[Service]
Type=oneshot
ExecStart=/opt/pdf/expansion/integrity_chain.py
