# infra/proxmox-export-2026-03-14/ct123/pdf-historian-snapshot.service
[Unit]
Description=PDF Historian Snapshot Writer
After=pdf-mission-control.service
Wants=pdf-mission-control.service

[Service]
Type=oneshot
ExecStart=/opt/pdf/mission-control/historian_snapshot.py
