# infra/proxmox-export-2026-03-14/ct123/pdf-historian-snapshot.timer
[Unit]
Description=Run PDF Historian snapshot every 15 min

[Timer]
OnBootSec=4min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
