# infra/proxmox-export-2026-03-14/ct123/pdf-night-watch.service
[Unit]
Description=PDF Night Watch urgent escalation
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
Environment=WEBHOOK_URL=REDACTED_WEBHOOK_URL
ExecStart=/opt/pdf/next3/night_watch.py
